# gost-watcher.path
[Path]
PathModified=/usr/local/etc/gost/

[Install]
WantedBy=multi-user.target
